# dist/debian/systemd/scylla-server.service.d/debian.conf
[Service]
EnvironmentFile=
EnvironmentFile=/etc/default/scylla-server
EnvironmentFile=/etc/scylla.d/*.conf
