# dist/common/systemd/scylla-fstrim.service
[Unit]
Description=Scylla fstrim
After=network.target

[Service]
Type=simple
ExecStart=/opt/scylladb/lib/scylla/scylla_fstrim

[Install]
WantedBy=multi-user.target
