# dist/debian/scylla-server.cron.d
0 0 * * sat   root	/opt/scylladb/lib/scylla/scylla_delay_fstrim
